# meta-common/network/usb-network/usb-network.service
[Unit]
Description=Enable USB Network

[Service]
Type=oneshot
ExecStart=/usr/sbin/usb-network.sh

[Install]
WantedBy=multi-user.target
