# meta-runbmc-nuvoton/recipes-graphics/obmc-ikvm/obmc-ikvm/obmc-ikvm.service
[Unit]
Description=obmc-ikvm

[Service]
ExecStart={bindir}/obmc-ikvm

[Install]
RequiredBy={SYSTEMD_DEFAULT_TARGET}
